# meta-inventec/meta-starscream/recipes-phosphor/init/starscream-init/starscream-init.service
[Unit]
Description=Essential init commands for starscream
Before=xyz.openbmc_project.EntityManager.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/sbin/starscream-init.sh

[Install]
WantedBy=multi-user.target
